# .dotfiles/ansible/roles/greenclip/templates/greenclip.service.j2
[Unit]
Description=Start greenclip daemon
After=display-manager.service

[Service]
ExecStart={{ ansible_env.HOME }}/.local/bin/greenclip daemon
Restart=always

[Install]
WantedBy=default.target
